# kodi/upstart/kodi.conf
env USER=kodi
start on (filesystem and stopped udevtrigger)
stop on runlevel [016]
respawn
script
  exec su -c "xinit /usr/bin/kodi --standalone -- -nocursor :0" $USER
end script