# iptables/etc/network/if-pre-up.d/iptables
#!/bin/bash
/usr/local/bin/iptables-redo
