# xbmc/upstart/xbmc.conf
env USER=xbmc
start on (filesystem and stopped udevtrigger)
stop on runlevel [016]
respawn
script
  exec su -c "xinit /usr/bin/xbmc --standalone -- -nocursor :0" $USER
end script