# xbmc/90-remote-wake.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="0609", ATTRS{idProduct}=="0334" RUN+="/bin/sh -c 'echo enabled > /sys$env{DEVPATH}/../power/wakeup'"