# xbmc/ir-keytable.conf
start on runlevel [2345]
stop on runlevel [06]

script
/usr/bin/ir-keytable -c
/usr/bin/ir-keytable -c -p RC-6,LIRC -w /lib/udev/rc_keymaps/rc6_mce
end script
